# packaging/prerm
#!/bin/sh
set -e

# Stop and disable the DX-RT service before the package is removed
if [ -d /run/systemd/system ]; then
    systemctl disable --now dxrt.service || true
fi

exit 0
